# deploy/telegram-bot.service.example
[Unit]
Description=Telegram Bot for Appointment System (FastAPI)
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/path/to/appointment-system
EnvironmentFile=/path/to/appointment-system/.env
ExecStart=/path/to/appointment-system/venv/bin/python -m bot.run
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
